# deploy/blow-your-work.service
[Unit]
Description=Blow Your Work FastAPI service
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/blow-your-work
EnvironmentFile=/home/azureuser/blow-your-work/.env
ExecStart=/home/azureuser/blow-your-work/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
